# 50-cxi-driver.rules
KERNEL=="cxi_ss1", SUBSYSTEM=="module", ACTION=="add", RUN+="/sbin/modprobe -a cxi_user cxi_eth"
